# root/etc/s6/inotify/run
#!/bin/bash

s6-svwait -U /etc/s6/syncthing

echo "Syncing $VOL"

exec syncthing-inotify -target="http://127.0.0.1:8080"