# root/etc/s6/inotify/run
#!/bin/bash

s6-svwait -u /etc/s6/syncthing

echo "Syncing $VOL"

exec inotifywait -dm syncthing-inotify