# root/etc/s6/inotify/run
#!/bin/bash

s6-svwait -u syncthing

echo "Syncing $VOL"

exec inotify -dm syncthing-inotify