# root/etc/s6/inotify/run
#!/bin/bash

s6-svwait -u syncthing

exec inotify -dm syncthing-inotify