# root/etc/s6/syncthing/run
#!/bin/bash

exec syncthing -gui-address="0.0.0.0:8080" -no-browser