# root/etc/s6/syncthing/run
#!/bin/bash

echo s6-svok /etc/s6/consul

s6-svwait -u /etc/s6/consul

if [ "$VOL" == "any" ]
then
	VOL="default"
fi

exec syncthing -gui-address="0.0.0.0:8080" -no-browser