# root/etc/s6/configure/finish
#!/bin/bash

s6-svc -O /etc/s6/configure